# frontend/supervisord.conf
[supervisord]
nodaemon=true

[program:receiver]
command=/usr/local/bin/python3 /python/app/receiver.py

[program:frontend]
command=/usr/local/bin/python3 /python/app/entrypoint.py